# systemd/naga.service
[Unit]
Description=Naga

[Service]
ExecStart=/usr/local/bin/naga
