# systemd/xsession.target
[Unit]
Description = X session is up and running
BindsTo=graphical-session.target
